# files-organizer.service
[Unit]
Description="Organizes the files in Download folder"

[Service]
ExecStart=/usr/bin/python3 /usr/local/Files-Organizer/organizer.py
